# scripts/launchd/com.contextdna.ecosystem-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contextdna.ecosystem-health</string>

    <key>ProgramArguments</key>
    <array>
        <string>__USER_HOME__/dev/er-simulator-superrepo/.venv/bin/python3</string>
        <string>__USER_HOME__/dev/er-simulator-superrepo/memory/ecosystem_health.py</string>
        <string>--daemon</string>
        <string>--auto-recover</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/dev/er-simulator-superrepo/memory/.ecosystem_health_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/dev/er-simulator-superrepo/memory/.ecosystem_health_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__USER_HOME__/dev/er-simulator-superrepo</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__USER_HOME__/dev/er-simulator-superrepo</string>
    </dict>

    <!-- Start 45 seconds after login to allow Docker to start -->
    <key>ThrottleInterval</key>
    <integer>45</integer>

    <!-- Match unified plist FD limits for SQLite + HTTP + log access -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>
</dict>
</plist>
