# scripts/launchd/com.considerate.vllm-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.considerate.vllm-watchdog</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__USER_HOME__/dev/er-simulator-superrepo/scripts/vllm-watchdog.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>30</integer>  <!-- Check every 30 seconds -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/vllm-watchdog-layer2.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/vllm-watchdog-layer2-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
